# debian/jenkins-agent.service
[Unit]
Description=Simple demonstration service

[Service]
ExecStart=/usr/bin/jenkins-agent

[Install]
WantedBy=multi-user.target